# zeer-telegram.service
# Systemd service file for zeer Telegram bot
# 
# Installation:
# 1. Edit this file and replace placeholders:
#    - /path/to/zeer with your actual path
#    - your_username with your Linux username
# 2. Copy to systemd: sudo cp zeer-telegram.service /etc/systemd/system/
# 3. Reload systemd: sudo systemctl daemon-reload
# 4. Enable service: sudo systemctl enable zeer-telegram
# 5. Start service: sudo systemctl start zeer-telegram
# 6. Check status: sudo systemctl status zeer-telegram
# 7. View logs: sudo journalctl -u zeer-telegram -f

[Unit]
Description=zeer Telegram Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=your_username
Group=your_username
WorkingDirectory=/path/to/zeer
Environment="PATH=/path/to/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/path/to/venv/bin/python start_telegram_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/path/to/zeer

[Install]
WantedBy=multi-user.target
